# ports/extra/mesa/.footprint
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/include/
drwxr-xr-x	root/root	usr/include/EGL/
-rw-r--r--	root/root	usr/include/EGL/eglext_angle.h
-rw-r--r--	root/root	usr/include/EGL/eglmesaext.h
drwxr-xr-x	root/root	usr/include/GL/
drwxr-xr-x	root/root	usr/include/GL/internal/
-rw-r--r--	root/root	usr/include/GL/internal/dri_interface.h
-rw-r--r--	root/root	usr/include/gbm.h
-rw-r--r--	root/root	usr/include/xa_composite.h
-rw-r--r--	root/root	usr/include/xa_context.h
-rw-r--r--	root/root	usr/include/xa_tracker.h
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/dri/
-rwxr-xr-x	root/root	usr/lib/dri/crocus_dri.so
-rwxr-xr-x	root/root	usr/lib/dri/i915_dri.so
-rwxr-xr-x	root/root	usr/lib/dri/iris_dri.so
-rwxr-xr-x	root/root	usr/lib/dri/kms_swrast_dri.so
-rwxr-xr-x	root/root	usr/lib/dri/nouveau_dri.so
-rwxr-xr-x	root/root	usr/lib/dri/r300_dri.so
-rwxr-xr-x	root/root	usr/lib/dri/r600_dri.so
-rwxr-xr-x	root/root	usr/lib/dri/radeonsi_dri.so
-rwxr-xr-x	root/root	usr/lib/dri/swrast_dri.so
-rwxr-xr-x	root/root	usr/lib/dri/virtio_gpu_dri.so
-rwxr-xr-x	root/root	usr/lib/dri/vmwgfx_dri.so
lrwxrwxrwx	root/root	usr/lib/libEGL_mesa.so -> libEGL_mesa.so.0
lrwxrwxrwx	root/root	usr/lib/libEGL_mesa.so.0 -> libEGL_mesa.so.0.0.0
-rwxr-xr-x	root/root	usr/lib/libEGL_mesa.so.0.0.0
lrwxrwxrwx	root/root	usr/lib/libGLX_mesa.so -> libGLX_mesa.so.0
lrwxrwxrwx	root/root	usr/lib/libGLX_mesa.so.0 -> libGLX_mesa.so.0.0.0
-rwxr-xr-x	root/root	usr/lib/libGLX_mesa.so.0.0.0
lrwxrwxrwx	root/root	usr/lib/libgbm.so -> libgbm.so.1
lrwxrwxrwx	root/root	usr/lib/libgbm.so.1 -> libgbm.so.1.0.0
-rwxr-xr-x	root/root	usr/lib/libgbm.so.1.0.0
lrwxrwxrwx	root/root	usr/lib/libglapi.so -> libglapi.so.0
lrwxrwxrwx	root/root	usr/lib/libglapi.so.0 -> libglapi.so.0.0.0
-rwxr-xr-x	root/root	usr/lib/libglapi.so.0.0.0
lrwxrwxrwx	root/root	usr/lib/libxatracker.so -> libxatracker.so.2
lrwxrwxrwx	root/root	usr/lib/libxatracker.so.2 -> libxatracker.so.2.5.0
-rwxr-xr-x	root/root	usr/lib/libxatracker.so.2.5.0
drwxr-xr-x	root/root	usr/lib/pkgconfig/
-rw-r--r--	root/root	usr/lib/pkgconfig/dri.pc
-rw-r--r--	root/root	usr/lib/pkgconfig/gbm.pc
-rw-r--r--	root/root	usr/lib/pkgconfig/xatracker.pc
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/drirc.d/
-rw-r--r--	root/root	usr/share/drirc.d/00-mesa-defaults.conf
drwxr-xr-x	root/root	usr/share/glvnd/
drwxr-xr-x	root/root	usr/share/glvnd/egl_vendor.d/
-rw-r--r--	root/root	usr/share/glvnd/egl_vendor.d/50_mesa.json
