# templates/mesa-32/.files
drwxr-xr-x root/root usr/
drwxr-xr-x root/root usr/share/
drwxr-xr-x root/root usr/include/
drwxr-xr-x root/root usr/lib/
-rwxr-xr-x root/root usr/lib/libGLESv2.so.2.0.0
lrwxrwxrwx root/root usr/lib/libglapi.so
lrwxrwxrwx root/root usr/lib/libxatracker.so.2
lrwxrwxrwx root/root usr/lib/libgbm.so
lrwxrwxrwx root/root usr/lib/libgbm.so.1
lrwxrwxrwx root/root usr/lib/libGLESv2.so
lrwxrwxrwx root/root usr/lib/libGL.so
lrwxrwxrwx root/root usr/lib/libGLESv1_CM.so
lrwxrwxrwx root/root usr/lib/libEGL.so.1
lrwxrwxrwx root/root usr/lib/libGLESv1_CM.so.1
lrwxrwxrwx root/root usr/lib/libGLESv2.so.2
-rwxr-xr-x root/root usr/lib/libgbm.so.1.0.0
lrwxrwxrwx root/root usr/lib/libGL.so.1
-rwxr-xr-x root/root usr/lib/libxatracker.so.2.5.0
drwxr-xr-x root/root usr/lib/dri/
lrwxrwxrwx root/root usr/lib/libEGL.so
-rwxr-xr-x root/root usr/lib/libEGL.so.1.0.0
-rwxr-xr-x root/root usr/lib/libGL.so.1.2.0
lrwxrwxrwx root/root usr/lib/libglapi.so.0
drwxr-xr-x root/root usr/lib/vdpau/
drwxr-xr-x root/root usr/lib/pkgconfig/
-rwxr-xr-x root/root usr/lib/libGLESv1_CM.so.1.1.0
-rwxr-xr-x root/root usr/lib/libglapi.so.0.0.0
lrwxrwxrwx root/root usr/lib/libxatracker.so
-rw-r--r-- root/root usr/lib/pkgconfig/glesv2.pc
-rw-r--r-- root/root usr/lib/pkgconfig/gbm.pc
-rw-r--r-- root/root usr/lib/pkgconfig/xatracker.pc
-rw-r--r-- root/root usr/lib/pkgconfig/glesv1_cm.pc
-rw-r--r-- root/root usr/lib/pkgconfig/egl.pc
-rw-r--r-- root/root usr/lib/pkgconfig/gl.pc
-rw-r--r-- root/root usr/lib/pkgconfig/dri.pc
lrwxrwxrwx root/root usr/lib/vdpau/libvdpau_r600.so
-rwxr-xr-x root/root usr/lib/vdpau/libvdpau_virtio_gpu.so.1.0.0
hrwxr-xr-x root/root usr/lib/vdpau/libvdpau_r600.so.1.0.0
lrwxrwxrwx root/root usr/lib/vdpau/libvdpau_radeonsi.so.1
lrwxrwxrwx root/root usr/lib/vdpau/libvdpau_radeonsi.so.1.0
lrwxrwxrwx root/root usr/lib/vdpau/libvdpau_virtio_gpu.so.1.0
hrwxr-xr-x root/root usr/lib/vdpau/libvdpau_nouveau.so.1.0.0
lrwxrwxrwx root/root usr/lib/vdpau/libvdpau_r600.so.1
lrwxrwxrwx root/root usr/lib/vdpau/libvdpau_nouveau.so.1
lrwxrwxrwx root/root usr/lib/vdpau/libvdpau_virtio_gpu.so.1
lrwxrwxrwx root/root usr/lib/vdpau/libvdpau_r600.so.1.0
lrwxrwxrwx root/root usr/lib/vdpau/libvdpau_nouveau.so.1.0
lrwxrwxrwx root/root usr/lib/vdpau/libvdpau_radeonsi.so
lrwxrwxrwx root/root usr/lib/vdpau/libvdpau_virtio_gpu.so
lrwxrwxrwx root/root usr/lib/vdpau/libvdpau_nouveau.so
hrwxr-xr-x root/root usr/lib/vdpau/libvdpau_radeonsi.so.1.0.0
-rwxr-xr-x root/root usr/lib/dri/i915_dri.so
hrwxr-xr-x root/root usr/lib/dri/iris_dri.so
-rwxr-xr-x root/root usr/lib/dri/virtio_gpu_drv_video.so
hrwxr-xr-x root/root usr/lib/dri/radeonsi_dri.so
hrwxr-xr-x root/root usr/lib/dri/nouveau_dri.so
hrwxr-xr-x root/root usr/lib/dri/vmwgfx_dri.so
hrwxr-xr-x root/root usr/lib/dri/r600_dri.so
hrwxr-xr-x root/root usr/lib/dri/swrast_dri.so
hrwxr-xr-x root/root usr/lib/dri/nouveau_drv_video.so
hrwxr-xr-x root/root usr/lib/dri/kms_swrast_dri.so
hrwxr-xr-x root/root usr/lib/dri/r600_drv_video.so
hrwxr-xr-x root/root usr/lib/dri/radeonsi_drv_video.so
hrwxr-xr-x root/root usr/lib/dri/virtio_gpu_dri.so
hrwxr-xr-x root/root usr/lib/dri/crocus_dri.so
hrwxr-xr-x root/root usr/lib/dri/r300_dri.so
drwxr-xr-x root/root usr/include/EGL/
-rw-r--r-- root/root usr/include/gbm.h
drwxr-xr-x root/root usr/include/GL/
drwxr-xr-x root/root usr/include/GLES3/
drwxr-xr-x root/root usr/include/GLES2/
-rw-r--r-- root/root usr/include/xa_tracker.h
drwxr-xr-x root/root usr/include/KHR/
-rw-r--r-- root/root usr/include/xa_context.h
-rw-r--r-- root/root usr/include/xa_composite.h
drwxr-xr-x root/root usr/include/GLES/
-rw-r--r-- root/root usr/include/GLES/glext.h
-rw-r--r-- root/root usr/include/GLES/glplatform.h
-rw-r--r-- root/root usr/include/GLES/egl.h
-rw-r--r-- root/root usr/include/GLES/gl.h
-rw-r--r-- root/root usr/include/KHR/khrplatform.h
-rw-r--r-- root/root usr/include/GLES2/gl2ext.h
-rw-r--r-- root/root usr/include/GLES2/gl2platform.h
-rw-r--r-- root/root usr/include/GLES2/gl2.h
-rw-r--r-- root/root usr/include/GLES3/gl32.h
-rw-r--r-- root/root usr/include/GLES3/gl3ext.h
-rw-r--r-- root/root usr/include/GLES3/gl3.h
-rw-r--r-- root/root usr/include/GLES3/gl3platform.h
-rw-r--r-- root/root usr/include/GLES3/gl31.h
-rw-r--r-- root/root usr/include/GL/glxext.h
drwxr-xr-x root/root usr/include/GL/internal/
-rw-r--r-- root/root usr/include/GL/glcorearb.h
-rw-r--r-- root/root usr/include/GL/glx.h
-rw-r--r-- root/root usr/include/GL/glext.h
-rw-r--r-- root/root usr/include/GL/gl.h
-rw-r--r-- root/root usr/include/GL/internal/dri_interface.h
-rw-r--r-- root/root usr/include/EGL/eglmesaext.h
-rw-r--r-- root/root usr/include/EGL/eglext_angle.h
-rw-r--r-- root/root usr/include/EGL/eglext.h
-rw-r--r-- root/root usr/include/EGL/egl.h
-rw-r--r-- root/root usr/include/EGL/eglplatform.h
drwxr-xr-x root/root usr/share/drirc.d/
-rw-r--r-- root/root usr/share/drirc.d/00-mesa-defaults.conf
